# stow/bundle-desktop/.config/autostart/cosmic-disable-display.desktop
[Desktop Entry]
Name=Cosmic disable eDP-1
Exec=sh -c "~/.local/scripts/cosmic-disable-display"
Terminal=false
Type=Application
